# allinone/scripts/all-in-one.conf
[supervisord]
nodaemon=true
loglevel=debug
logfile=/usr/local/supervisord.log

[program:zookeeper]
user=root
command=/usr/local/zookeeper/bin/zkServer.sh start-foreground
stderr_logfile=/usr/local/zookeeper/zookeeper-supervisor.log
priority=100

[program:mysql]
command=/usr/local/scripts/start-mysql.sh
redirect_stderr=true
priority=100

[program:init-db]
user=root
command=/usr/local/scripts/wait-for-it.sh localhost:3306 -- /usr/local/scripts/init-db.sh
priority=200

[program:druid-coordinator]
directory=/usr/local/druid
command=java
  -server
  -Xms1g
  -Xmx3g
  -XX:MaxDirectMemorySize=2g
  -verbose:gc
  -XX:+PrintGCDetails
  -XX:+PrintGCTimeStamps
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Djava.io.tmpdir=%(ENV_DRUID_HOME)s/tmp
  -Djava.util.logging.manager=org.apache.logging.log4j.jul.LogManager
  -Ddruid.extensions.directory=%(ENV_DRUID_HOME)s/extensions
  -Ddruid.coordinator.asOverlord.enabled=true
  -Ddruid.coordinator.asOverlord.overlordService=druid/overlord
  -Ddruid.indexer.fork.property.druid.processing.numThreads=1
  -Ddruid.indexer.storage.type=metadata
  -Ddruid.indexer.queue.startDelay=PT0M
  -Ddruid.indexer.runner.javaOpts="-server -Xmx3g -XX:MaxDirectMemorySize=2147483648"
  -Ddruid.processing.buffer.sizeBytes=536870912
  -Ddruid.coordinator.startDelay=PT5S
  -cp %(ENV_DRUID_HOME)s/conf/druid:%(ENV_DRUID_HOME)s/lib/*
  io.druid.cli.Main server coordinator
stderr_logfile = /usr/local/druid/logs/druid-coordinator-stderr.log
stdout_logfile = /usr/local/druid/logs/druid-coordinator-stdout.log
user=root
priority=300

[program:druid-historical]
directory=/usr/local/druid
command=java
  -server
  -Xms16g
  -Xmx16g
  -XX:MaxDirectMemorySize=8g
  -verbose:gc
  -XX:+PrintGCDetails
  -XX:+PrintGCTimeStamps
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Djava.io.tmpdir=%(ENV_DRUID_HOME)s/tmp
  -Djava.util.logging.manager=org.apache.logging.log4j.jul.LogManager
  -Ddruid.extensions.directory=%(ENV_DRUID_HOME)s/extensions
  -Ddruid.computation.buffer.size=67108864
  -Ddruid.server.maxSize=5000000000
  -cp %(ENV_DRUID_HOME)s/conf/druid:%(ENV_DRUID_HOME)s/lib/*
  io.druid.cli.Main server historical
stderr_logfile = /usr/local/druid/logs/druid-historical-stderr.log
stdout_logfile = /usr/local/druid/logs/druid-historical-stdout.log
user=root
priority=300

[program:druid-broker]
directory=/usr/local/druid
command=java
  -server
  -Xms24g
  -Xmx24g
  -XX:MaxDirectMemorySize=8g
  -verbose:gc
  -XX:+PrintGCDetails
  -XX:+PrintGCTimeStamps
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Djava.io.tmpdir=%(ENV_DRUID_HOME)s/tmp
  -Djava.util.logging.manager=org.apache.logging.log4j.jul.LogManager
  -Ddruid.computation.buffer.size=67108864
  -Ddruid.broker.cache.sizeInBytes=33554432
  -Ddruid.extensions.directory=%(ENV_DRUID_HOME)s/extensions
  -cp %(ENV_DRUID_HOME)s/conf/druid:%(ENV_DRUID_HOME)s/lib/*
  io.druid.cli.Main server broker
stderr_logfile = /usr/local/druid/logs/druid-broker-stderr.log
stdout_logfile = /usr/local/druid/logs/druid-broker-stdout.log
user=root
priority=300

[program:druid-overlord]
directory=/usr/local/druid
command=java
  -server
  -Xms3g
  -Xmx3g
  -verbose:gc
  -XX:+PrintGCDetails
  -XX:+PrintGCTimeStamps
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Djava.io.tmpdir=%(ENV_DRUID_HOME)s/tmp
  -Djava.util.logging.manager=org.apache.logging.log4j.jul.LogManager
  -Ddruid.computation.buffer.size=67108864
  -Ddruid.extensions.directory=%(ENV_DRUID_HOME)s/extensions
  -Ddruid.indexer.queue.startDelay=PT5S
  -Ddruid.indexer.runner.javaOpts=-server -Xmx1g
  -Ddruid.indexer.fork.property.druid.processing.numThreads=2
  -Ddruid.indexer.fork.property.druid.computation.buffer.size=100000000
  -Ddruid.indexer.fork.property.druid.indexer.task.defaultHadoopCoordinates="[\"org.apache.hadoop:hadoop-client:%(ENV_HADOOP_VERSION)s\"]"
  -Ddruid.indexer.fork.property.druid.server.http.numThreads=50
  -Ddruid.indexer.storage.type=metadata
  -cp %(ENV_DRUID_HOME)s/conf/druid:%(ENV_DRUID_HOME)s/lib/*
  io.druid.cli.Main server overlord
stderr_logfile = /usr/local/druid/logs/druid-overlord-stderr.log
stdout_logfile = /usr/local/druid/logs/druid-overlord-stdout.log
user=root
priority=300

[program:druid-middleManager]
directory=/usr/local/druid
command=java
  -server
  -Xms4g
  -Xmx4g
  -verbose:gc
  -XX:+PrintGCDetails
  -XX:+PrintGCTimeStamps
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Djava.io.tmpdir=%(ENV_DRUID_HOME)s/tmp
  -Djava.util.logging.manager=org.apache.logging.log4j.jul.LogManager
  -Ddruid.computation.buffer.size=67108864
  -Ddruid.extensions.directory=%(ENV_DRUID_HOME)s/extensions
  -cp %(ENV_DRUID_HOME)s/conf/druid:%(ENV_DRUID_HOME)s/lib/*
  io.druid.cli.Main server middleManager
stderr_logfile = /usr/local/druid/logs/druid-middlemanager-stderr.log
stdout_logfile = /usr/local/druid/logs/druid-middlemanager-stdout.log
user=root
priority=300


[program:druid-router]
directory=/usr/local/druid
command=java
  -server
  -Xms512m
  -Xmx512m
  -XX:+UseG1GC
  -XX:MaxDirectMemorySize=512m
  -XX:+ExitOnOutOfMemoryError
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Djava.io.tmpdir=var/tmp
  -Djava.util.logging.manager=org.apache.logging.log4j.jul.LogManager
  -Ddruid.service=druid/router
  -Ddruid.plaintextPort=8888
  -Ddruid.router.http.numConnections=50
  -Ddruid.router.http.readTimeout=PT5M
  -Ddruid.router.http.numMaxThreads=100
  -Ddruid.server.http.numThreads=100
  -Ddruid.router.defaultBrokerServiceName=druid/broker
  -Ddruid.router.coordinatorServiceName=druid/coordinator
  -Ddruid.router.managementProxy.enabled=true
  -cp %(ENV_DRUID_HOME)s/conf/druid:%(ENV_DRUID_HOME)s/lib/*
  io.druid.cli.Main server router
stderr_logfile = /usr/local/druid/logs/druid-router-stderr.log
stdout_logfile = /usr/local/druid/logs/druid-router-stdout.log
user=root
priority=300

[program:discovery]
user=root
command=bash -c 'if [ ${PERSIST_DISCOVERY} = "TRUE" ];then /usr/local/scripts/wait-for-it.sh localhost:8082 -- /usr/local/metatron-discovery/bin/metatron.sh start;
  else /usr/local/scripts/wait-for-it.sh localhost:8082 -- /usr/local/metatron-discovery/bin/metatron.sh --init start;fi'
redirect_stderr=true
priority=400
